# udev/99-local-devcoredump.rules.in
ACTION=="add", SUBSYSTEM=="devcoredump", RUN+="@sbindir@/devcoredump /sys/%p"
